# .config/systemd/user/snapper-weekly@.timer
[Unit]
Description=Weekly Snapper snapshot timer for %i

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
